# debian/extra/systemd/svx-history-rotate.service
[Unit]
Description=Rotate Svx dashboard history
After=svx-sse.service

[Service]
Type=oneshot
ExecStart=/usr/sbin/svx-history-rotate
